# manifests/syncthing.rc.erb
#!/bin/sh
#
# PROVIDE: syncthing
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name=syncthing
rcvar=syncthing_enable

command="/usr/local/sbin/utility"

load_rc_config $name

#
# DO NOT CHANGE THESE DEFAULT VALUES HERE
# SET THEM IN THE /etc/rc.conf FILE
#
utility_enable=${utility_enable-"NO"}
pidfile=${utility_pidfile-"/var/run/utility.pid"}

run_rc_command "$1"
